# source/services/acpid.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=ACPI daemon

[Service]
ExecStart=%p -f -d -l
StandardInput=socket

[Install]
WantedBy=workstation.target
