# source/services/amd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=auto-mount daemon
Requires=rpcbind.service nfsclient.service
After=rpcbind.service nfsclient.service

[Service]
EnvironmentDirectory=env
ExecStart=%p ${flags}

[Install]
WantedBy=workstation.target
