# source/services/anacron.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=anacron scheduling service

[Service]
ExecStartPre=on_ac_power
ExecStart=%p -d -s
Restart=never

[Install]
WantedBy=workstation.target
