# source/services/apache2.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Apache HTTP server

[Service]
EnvironmentDirectory=env
#User=%p
RuntimeDirectory=%p
RuntimeDirectoryMode=0755
LimitNOFILE=8192
ExecStart=%p ${root:+-d "${root}"} ${conffile:+-f "${conffile}"} -k start
PrivateTmp=true

[Install]
WantedBy=server.target
