# source/services/avahi-daemon.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Bonjour service
After=dbus.service

[Service]
systemdWorkingDirectory=false
ExecStartPre=mkdir -m 02755 /run/avahi-daemon
ExecStartPre=chown avahi /run/avahi-daemon
ExecStart=%P
Restart=always

[Install]
WantedBy=workstation.target
