# source/services/background-fsck.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Background filesystem check (deprecated in favour of fsck@device service bundles)
DefaultDependencies=false

[Service]
EnvironmentDirectory=env
ExecStartPre=sleep ${delay:-0}
ExecStart=fsck -B -p

[Install]
StoppedBy=shutdown.target
WantedBy=sysinit.target
