# source/services/bcron-update.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Guenter's bcron update daemon

[Service]
systemdWorkingDirectory=false
EnvironmentDirectory=env
ExecStart=%P /etc/crontab /etc/cron.d
Restart=always

[Install]
WantedBy=server.target

