# source/services/bluetooth.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Bluetooth daemon
After=local-syslog.target

[Service]
Type=dbus
ExecStart=%pd -n
ProtectHome=true
LimitNPROC=1

[Install]
WantedBy=workstation.target
