# source/services/brltty-log.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Combined logger for brltty@ services

[Service]
User=%p
WorkingDirectory=/var/log/sv
ExecStart=cyclog brltty/
RestartSec=1

[Install]
WantedBy=workstation.target
