# source/services/bsd-serial-init.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=BSD serial device initialization service

[Service]
Type=oneshot
ExecStart=/etc/rc.d/serial
Restart=never
RemainAfterExit=yes

[Install]
WantedBy=workstation.target
