# source/services/bsnmpd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=basic SNMP service

[Service]
EnvironmentDirectory=env
ExecStart=%p -d

[Install]
WantedBy=server.target
