# source/services/ccd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=CCD service

[Service]
systemdWorkingDirectory=false
ExecStart=ccdconfig -C
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
