# source/services/ccd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=configure the Concatenated Disc Driver

[Service]
systemdWorkingDirectory=false
ExecStart=ccdconfig -C -v
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
