# source/services/cleanuucp.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

# Note that many BSDs no longer include UUCP.

[Unit]
Description=Clean-up UUCP spool directories

[Service]
Type=oneshot
WorkingDirectory=/var/spool/
ExecStart=sh -c 'rm -r uucp/.Temp/* ; rm -r lock/*'
RemainAfterExit=yes

[Install]
WantedBy=workstation.target
