# source/services/cnid_metad.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=AppleTalk CNID daemon spawning daemon
After=atalkd.service

[Service]
ExecStart=%p -d
Restart=always

[Install]
WantedBy=server.target
