# source/services/colord.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=ColourManager daemon
After=local-syslog.target

[Service]
Type=dbus
ExecStart=/usr/lib/x86_64-linux-gnu/%p/%p
Restart=always

[Install]
WantedBy=workstation.target
