# source/services/console-kit-daemon.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=ConsoleKit daemon

[Service]
Type=dbus
ExecStart=%p --no-daemon

[Install]
WantedBy=workstation.target
