# source/services/crashinfo.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Analyze the last saved core dump

[Service]
Type=oneshot
EnvironmentDirectory=env
ExecStart=%p ${flags} ${dumpdir:-/var/crash}
RemainAfterExit=yes

[Install]
WantedBy=workstation.target
