# source/services/cron.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Regular background program processing daemon

[Service]
ExecStart=%p -f
Restart=always

[Install]
WantedBy=workstation.target
