# source/services/ctld.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=CAM Target Layer/iSCSI target daemon
Wants=kmod@ctl
After=kmod@ctl

[Service]
# ctld needs a proper option for telling it not to fork.  Until it has, this will have to do.
ExecStart=%p -d

[Install]
WantedBy=workstation.target
