# source/services/cups-browsed.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=compatibility non-DNS-SD Common Unix Printing Service location
After=org.cups.cupsd.service avahi-daemon.service

[Service]
ExecStart=%p

[Install]
WantedBy=server.target
