# source/services/cyclog@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Standard format %P logging service for %I
Before=%i.service

[Service]
systemdUserEnvironment=no
WorkingDirectory=/var/log/sv/
User=%i-log
ExecStart=%P %I/
Restart=always

[Install]
WantedBy=workstation.target
