# source/services/ddclient.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Dynamic DNS update client

[Service]
User=%p
ExecStart=%p --foreground

[Install]
WantedBy=workstation.target
