# source/services/defaultroute.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=BSD wait for default route over DHCP.
Before=networking.target

[Service]
Type=oneshot
EnvironmentDirectory=env
ExecStart=./helper
RemainAfterExit=yes

[Install]
WantedBy=server.target
