# source/services/deluged.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Deluge BitTorrent client

[Service]
# FIXME: Is this really needed?
systemdUserEnvironment=yes
User=deluge
Umask=037
ExecStart=%p --do-not-daemonize

[Install]
WantedBy=server.target
