# source/services/dhclient@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=%p service for device %i

[Service]
systemdWorkingDirectory=false
ExecStart=%p -d %i
Restart=always

[Install]
WantedBy=workstation.target
