# source/services/dovecot.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Dovecot service

[Service]
ExecStart=%p -F
ProtectSystem=full

[Install]
WantedBy=server.target
