# source/services/epmd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Erlang Port Mapper service

[Service]
# FIXME: Is this really needed?
systemdUserEnvironment=true
User=rabbitmq
ExecStart=%p
Restart=always

[Install]
WantedBy=workstation.target
