# source/services/exim4-queue.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=exim4 queue runner daemon

[Service]
EnvironmentDirectory=env
LimitMEMORY=70000000
ExecStart=exim4 -q30m -d
ExecStartPre=exim4 -bV

[Install]
WantedBy=server.target
