# source/services/fancontrol.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=fan control daemon

[Service]
ExecStart=%p

[Install]
WantedBy=workstation.target
