# source/services/fcgiwrap.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=nginx's Fast CGI wrapper server

[Service]
systemdWorkingDirectory=false
User=http
ExecStart=%P
PrivateTmp=true
Restart=always

[Install]
WantedBy=server.target

