# source/services/fcgiwrap.socket
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=nginx's Fast CGI wrapper socket

[Socket]
Backlog=2
#MaxConnections=2
ListenStream=%t/%p.socket

[Install]
WantedBy=server.target
