# source/services/firstboot-pcbsd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Run PC-BSD first-boot task

[Service]
Type=oneshot
ExecStart=./helper
RemainAfterExit=yes

[Install]
WantedBy=workstation.target
