# source/services/freshclam.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Update ClamAV database daemon

[Service]
ExecStart=%p --daemon --stdout

[Install]
WantedBy=workstation.target
