# source/services/gpsd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=GPS daemon

[Service]
ExecStart=%p -N
StandardInput=socket

[Install]
