# source/services/hastd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=High Availability Storage service
Wants=kmod@geom_gate
After=kmod@geom_gate

[Service]
EnvironmentDirectory=env
# The daemon should properly clean itself up at SIGTERM; if it does not, that is a bug.
ExecStart=%p -F

[Install]
WantedBy=server.target
