# source/services/hostname.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Set the dynamic hostname
DefaultDependencies=false
Wants=local-fs.target
After=local-fs.target

[Service]
ExecStart=set-dynamic-hostname
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
