# source/services/hostapd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=IEEE 802.11 authenticator daemon

[Service]
EnvironmentDirectory=env
ExecStart=%p

[Install]
WantedBy=workstation.target
