# source/services/iocage.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=iocage start-up and shut-down "boot" jails
Before=securelevel.service

[Service]
ExecStart=%p rcboot
ExecStopPost=%p rcshutdown

[Install]
WantedBy=workstation.target
