# source/services/ipropd_slave.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Propagate Heimdal Kerberos KDC master to KDC slaves, slave side
After=kdc.service
Wants=kdc.service

[Service]
EnvironmentDirectory=env
ExecStart=%p --keytab "${keytab}" "${master}"

[Install]
WantedBy=server.target
