# source/services/iscsictl.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=iSCSI initiator control

[Service]
systemdWorkingDirectory=false
# The import procedure is expected to have populated this directory with flags.
EnvironmentDirectory=env
ExecStart=%p ${flags}
RemainAfterExit=true

[Install]
WantedBy=server.target
