# source/services/kdc.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=BSD Kerberos version 5 Authentication Service and Key Distribution Center

[Service]
EnvironmentDirectory=env
ExecStart=%p ${flags}

[Install]
WantedBy=server.target
