# source/services/kdm@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=KDE display manager using virtual terminal %i
After=hald.service genkdmconf.service

[Service]
EnvironmentDirectory=env
ExecStart=%p -nodaemon %i

[Install]
# TODO: This might be better as graphical.
WantedBy=workstation.target
