# source/services/keyserv.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=SunRPC keyserv daemon
Wants=rpcbind.service
After=rpcbind.service ypset.service

[Service]
EnvironmentDirectory=env
ExecStart=%p

[Install]
WantedBy=server.target
