# source/services/kfd@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=KF service

[Service]
systemdWorkingDirectory=false
ExecStart=%p
Restart=always

[Install]
WantedBy=server.target
