# source/services/klogd.socket
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]

[Socket]

[Install]
WantedBy=workstation.target
