# source/services/kmod@autofs.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=kmod loader for autofs4

[Service]
ExecStart=sh -c '[ -d /sys/class/misc/autofs4/ ] || system-control load-kernel-module autofs4'
RemainAfterExit=yes

[Install]
WantedBy=workstation.target
