# source/services/kpasswdd@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Kerberos v5 password daemon

[Service]
systemdWorkingDirectory=false
ExecStart=/usr/libexec/%p
Restart=always

[Install]
WantedBy=server.target
