# source/services/kpropd@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Kerberos propd slave

[Service]
ExecStart=%p

[Install]
WantedBy=server.target
