# source/services/lircd.socket
[Socket]
ListenStream=/run/lirc/lircd

[Install]
WantedBy=sockets.target
