# source/services/lpd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=line printer service
After=chkprintcap

[Service]
EnvironmentDirectory=env
ExecStart=%p

[Install]
WantedBy=server.target
