# source/services/machine-id.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Set up the machine ID
DefaultDependencies=false
After=local-fs.target

[Service]
ExecStart=setup-machine-id
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
