# source/services/mdadm.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Software RAID monitor daemon

[Service]
ExecStart=%p --monitor --scan

[Install]
WantedBy=workstation.target
