# source/services/miredo.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Teredo IP6 tunnelling client/relay

[Service]
ExecStart=%p --foreground

[Install]
WantedBy=workstation.target
