# source/services/monit.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=monitoring service

[Service]
ExecStart=%p -I
Restart=always

[Install]
WantedBy=workstation.target
