# source/services/monitor-fsck-progress.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=fsck combined progress information displayed on /dev/console
DefaultDependencies=off

[Service]
StandardInput=tty
StandardOutput=inherit
User=daemon
ExecStart=monitor-fsck-progress
TTYPath=console

[Install]
