# source/services/monitor-fsck-progress.socket
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=local socket used for %p

[Socket]
Backlog=2
SocketMode=0644
ListenStream=%t/fsck.progress
Accept=false
