# source/services/monitor-fsck-progress.socket
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=local socket used for %p
DefaultDependencies=off
Before=local-fs-pre.target

[Socket]
Backlog=2
ListenStream=/run/fsck.progress
Accept=false

[Install]
WantedBy=sysinit.target
