# source/services/mroute4d.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=IPv4 mroute daemon

[Service]
ExecStart=mrouted
Restart=always

[Install]
WantedBy=server.target
