# source/services/mroute6d.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=IPv6 mroute daemon

[Service]
EnvironmentDirectory=env
ExecStart=pim6dd

[Install]
WantedBy=server.target
