# source/services/multilog@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Standard format %p logging service for %I
Before=%I.service

[Service]
systemdUserEnvironment=no
WorkingDirectory=/var/log/sv/
User=%i-log
ExecStart=%p t ./%I/
Restart=always

[Install]
WantedBy=workstation.target
