# source/services/nagios-nrpe-server@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Nagios NRPE server

[Service]
User=nagios
ExecStart=nrpe -c /etc/nagios/nrpe.cfg --inetd

[Install]
WantedBy=server.target
