# source/services/nagios-nrpe-server@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Nagios NRPE server

[Service]
systemdWorkingDirectory=false
User=nagios
ExecStart=nrpe -c /etc/nagios/nrpe.cfg --inetd
Restart=no
StandardInput=socket

[Install]
WantedBy=server.target
