# source/services/nfscbd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=NFS client side callback daemon
After=nfsuserd.service

[Service]
EnvironmentDirectory=env
ExecStart=%p

[Install]
WantedBy=server.target
