# source/services/nisdomain.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=NIS domain name

[Service]
# The import procedure is expected to have populated this directory with name.
EnvironmentDirectory=env
ExecStart=domainname ${name}
RemainAfterExit=true

[Install]
WantedBy=workstation.target
