# source/services/ntpd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=NTP daemon

[Service]
EnvironmentDirectory=env
ExecStart=%p --nofork ${panicgate} ${user:+-u "${user}${group:+:${group}}"} ${logfile:+--logfile "${logfile}"} ${flags}

[Install]
WantedBy=workstation.target
