# source/services/openbsd-cron.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Regular background program processing daemon by Paul Vixie

[Service]
EnvironmentDirectory=env
ExecStart=cron -n

[Install]
WantedBy=workstation.target
