# source/services/org.cups.cups-cancel.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=cancel leftover Common Unix Printing Service jobs
Before=org.cups.cupsd.service

[Service]
Type=oneshot
ExecStart=cancel -a
RemainAfterExit=yes

[Install]
WantedBy=workstation.target
