# source/services/org.cups.cupsd.socket
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Common Unix Printing Service socket

[Socket]
Backlog=2
ListenStream=%t/cups/cups.sock

[Install]
WantedBy=workstation.target
