# source/services/packagekit.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=%p daemon
After=dbus.service
# Don't want D-Bus in case the D-Bus daemon is not managed by service-manager.
#Wants=dbus.service

[Service]
systemdWorkingDirectory=false
ExecStart=/usr/lib/%p/%pd
Restart=always

[Install]
WantedBy=workstation.target
