# source/services/pdnsd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Proxy DNS server

[Service]
ExecStart=%p

[Install]
WantedBy=workstation.target
