# source/services/pop3.socket
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Generic POP3 socket

[Socket]
Backlog=2
BindIPv6Only=both
ListenStream=pop3
Accept=true
MaxConnections=16
UCSPIRules=true
LogUCSPIRules=yes
NoDelay=false

[Install]
WantedBy=server.target
