# source/services/pulseaudio.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=PulseAudio Daemon
 
[Install]
WantedBy=multi-user.target
 
[Service]
# FIXME: Is this really needed?
systemdWorkingDirectory=on
PrivateTmp=true
ExecStart=pulseaudio --system --disallow-exit --disallow-module-loading
