# source/services/pwcheck.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=BSD password file lock check
DefaultDependencies=false

[Service]
ExecStart=sh -c '[ -f /etc/ptmp ] && echo 1>&2 "ERROR: password file may be incorrect -- /etc/ptmp exists"'
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
