# source/services/qmail-smtp-relay@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Bernstein's qmail SMTP Relay daemon

[Service]
EnvironmentDirectory=env
LimitMEMORY=70000000
User=qmaild
ExecStart=qmail-smtpd

[Install]
WantedBy=server.target
