# source/services/s6-fdholder.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=s6 fd holder daemon

[Service]
User=%p
ExecStart=%pd -i rules
StandardInput=socket

[Install]
