# source/services/saned@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Scanner service

[Service]
User=saned
ExecStart=%p
StandardInput=socket
Restart=always

[Install]
WantedBy=server.target
