# source/services/sendmail-smtp-relay@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Sendmail SMTP Relay daemon
# Sendmail abends invisibly if it cannot connect to a local syslog socket.
After=local-syslog.target

[Service]
EnvironmentDirectory=env
LimitMEMORY=70000000
ExecStart=sendmail -L %p -Am -bs ${flags}
ExecStartPre=sendmail -L %p -Am -bC

[Install]
WantedBy=server.target
