# source/services/sissylog@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=%p logging service for %I
Before=%I.service

[Service]
User=%i-log
ExecStart=%p %I

[Install]
WantedBy=workstation.target
