# source/services/slapd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=OpenLDAP %p daemon

[Service]
RuntimeDirectory=%p
RuntimeDirectoryMode=1777
EnvironmentDirectory=env
ExecStart=%p -d 0 ${SERVICES:+-h "${SERVICES}"} ${USER:+-u "${USER}"} ${GROUP:+-g "${GROUP}"} ${OPTIONS}

[Install]
WantedBy=server.target
