# source/services/slim.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Simple Lightweight Login Manager
After=ttylogin@ttyv6 ttylogin@tty7

[Service]
ExecStart=%p -nodaemon

[Install]
# TODO: This might be better as graphical.
WantedBy=workstation.target
