# source/services/slurpd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=OpenLDAP %p daemon

[Service]
ExecStart=%p -d 0

[Install]
WantedBy=server.target
