# source/services/smbd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=SAMBA file and print services daemon
After=nmbd.service

[Service]
ExecStart=%p -F -s "${conf:-/usr/local/etc/smb4.conf}"
LimitNOFILE=16384

[Install]
WantedBy=server.target
