# source/services/sshd@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=SSH per-connection service
Requires=sshdgenkeys.service
After=sshdgenkeys.service

[Service]
LocalReaper=true
EnvironmentDirectory=env
RuntimeDirectory=%p
RuntimeDirectoryMode=0755
ExecStart=%p -D -e -i

[Install]
WantedBy=server.target
