# source/services/sshd.socket
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=SSH socket

[Socket]
ListenStream=22
BindIPv6Only=both
Accept=true

[Install]
WantedBy=server.target
