# source/services/svnserve@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=svnserve per-connection service

[Service]
ExecStart=%p -i

[Install]
WantedBy=server.target
