# source/services/svnserve@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=svnserve per-connection service

[Service]
systemdWorkingDirectory=false
ExecStart=%p -i
Restart=always

[Install]
WantedBy=server.target
