# source/services/swift@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=openstack %p %i server

[Service]
systemdWorkingDirectory=false
EnvironmentFile=-/etc/sysconfig/%p-%i
ExecStart=%p-%i /etc/%p/%i-server.conf
Restart=always

[Install]
WantedBy=server.target
