# source/services/swift@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=OpenStack %p %i server

[Service]
#EnvironmentFile=-/etc/sysconfig/%p-%i
EnvironmentDirectory=env
ExecStart=%p-%i /etc/%p/%i-server.conf
Restart=always

[Install]
WantedBy=server.target
