# source/services/systemd-syslog.socket
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=syslog socket used when attaching syslog daemons to the systemd journal

[Socket]
ListenDatagram=%t/systemd/journal/syslog
SocketMode=0666
PassCredentials=yes
PassSecurity=yes

[Install]
WantedBy=workstation.target
