# source/services/tinydns.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Bernstein's content DNS server daemon

[Service]
EnvironmentUser=tinydns-d
EnvironmentDirectory=env
LimitDATA=300000
ExecStart=%p
Restart=always

[Install]
WantedBy=server.target
