# source/services/uhidd-log.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Combined logger for uhidd@ services

[Service]
User=%p
WorkingDirectory=/var/log/sv
ExecStart=cyclog uhidd/
RestartSec=1

[Install]
WantedBy=workstation.target
