# source/services/uhidd@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=%p service for device %I
Wants=kmod@vkbd
After=kmod@vkbd

[Service]
systemdWorkingDirectory=false
ExecStart=%p -d %I
Restart=always

[Install]
WantedBy=workstation.target
