# source/services/update-binfmts.service
#
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=additional executable binary formats
Conflicts=binfmt-support.service

[Service]
ExecStart=%p --enable
ExecStopPost=%p --disable
RemainAfterExit=true

[Install]
WantedBy=workstation.target
