# source/services/uwsgi-emperor.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=uWSGI "Emperor" listening server

[Service]
RuntimeDirectory=uwsgi
ExecStart=uwsgi --ini /etc/uwsgi/emperor.ini

[Install]
WantedBy=server.target
