# source/services/VBoxService.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=VirtualBox guest additions daemon

[Service]
ExecStart=VBoxService -f

[Install]
WantedBy=virtualbox-guest.target
