# source/services/virecover.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=nvi recovery

[Service]
Type=oneshot
EnvironmentDirectory=env
ExecStart=./helper
Restart=never

[Install]
WantedBy=workstation.target
