# source/services/vixiedhcp6d@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Vixie/ISC DHCP over IPv6 server on %i

[Service]
ExecStart=dhcpd -6 -f -d -q %i

[Install]
WantedBy=server.target
