# source/services/vsftpd-ssl@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=FTP service over IP4 using vsftp

[Service]
# FIXME: Is this really needed?
systemdUserEnvironment=true
LimitData=300000
User=vsftpd
ExecStart=%p -oimplicit_ssl=yes
Restart=always

[Install]
WantedBy=server.target
