# source/services/wpa_supplicant.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Wi-Fi Protected Access client daemon

[Service]
Type=dbus
ExecStart=%p -u
Restart=always

[Install]
WantedBy=workstation.target
