# source/services/ypbind.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=NIS binding daemon

[Service]
Type=forking
# The import procedure is expected to have populated this directory with flags.
EnvironmentDirectory=env
ExecStart=%p ${flags}

[Install]
WantedBy=server.target
