# source/convert/per-user/gnome-keyring-daemon.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=GNOME keyring (secrets only) for %u

[Service]
Type=dbus
EnvironmentDirectory=env
ExecStart=%p --foreground --start --components=secrets
Restart=on-abort

[Install]
