# source/convert/per-user/mpd.socket
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Music Player socket used by %u

[Socket]
ListenStream=%t/%p/socket
Backlog=5
KeepAlive=true
PassCredentials=true

[Install]
WantedBy=normal.target
