# source/convert/per-user/pulseaudio.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=PulseAudio daemon for %u
 
[Service]
RuntimeDirectory=pulse
EnvironmentDirectory=env
ExecStart=%p --exit-idle-time=-1 --daemonize=no ${flags}
Restart=on-abort
 
[Install]
WantedBy=normal.target
