# source/convert/per-user/startup.target
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Startup action group target
Before=shutdown
After=intrat
Wants=intrat

[Install]
StoppedBy=shutdown
