# source/convert/per-user/telepathy-mission-control-5.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Telepathy account manager for %u

[Service]
Type=dbus
EnvironmentDirectory=env
ExecStart=/usr/lib/telepathy/mission-control-5
Restart=on-abort

[Install]
