# source/convert/per-user/urxvt.socket
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=RXVT (Unicode) daemon control socket for %u

[Socket]
ListenStream=%t/%p/daemon
SocketMode=0600

[Install]
