# source/convert/per-user/watchman.socket
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Facebook Watchman daemon control socket for %u

[Socket]
ListenStream=%t/%p/state/sock
SocketMode=0600

[Install]
