# source/convert/per-user/zeitgeist-datahub.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Zeitgeist datahub for %u

[Service]
Type=dbus
EnvironmentDirectory=env
ExecStart=%p
Restart=on-abort

[Install]
