# source/convert/per-user/zeitgeist-daemon.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Zeitgeist activity monitor for %u

[Service]
Type=dbus
EnvironmentDirectory=env
ExecStart=%p
Restart=on-abort

[Install]
