# source/targets/jails.target
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Jails and jailed daemons collection

[Service]
User=nobody

[Install]
StoppedBy=shutdown.target
WantedBy=workstation.target
