# source/targets/local-fs-pre.target
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Target for things to complete before mount/fsck of local filesystems
EarlySupervise=true
Before=shutdown.target

[Service]
User=nobody

[Install]
StoppedBy=shutdown.target
WantedBy=sysinit.target
