# source/targets/multi-user-pre.target
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Target for things to complete before running login services
EarlySupervise=true

[Service]
User=nobody

[Install]
StoppedBy=shutdown.target
