# source/targets/multi-user.target
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Multi-user group target
EarlySupervise=true
Wants=multi-user-pre.target

[Service]
User=nobody

[Install]
StoppedBy=shutdown.target
WantedBy=normal.target
