# source/targets/sockets.target
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Sockets milestone for systemd compatibility

[Service]
User=nobody

[Install]
StoppedBy=shutdown.target
