# source/targets/uhidd.target
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=%p umbrella

[Service]
User=nobody

[Install]
WantedBy=workstation.target
