# source/targets/unmount.target
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Unmount all mount points action target
EarlySupervise=true
After=shutdown.target

[Service]
User=nobody

[Install]
