# source/targets/virtualbox-guest.target
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=VirtualBox guest machine target
EarlySupervise=true

[Service]
User=nobody

[Install]
WantedBy=workstation.target
