# source/systemd/service-manager-svscan.path
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Daemontools compatible /service/ directory
After=local-fs.target sysinit.target

[Path]
DirectoryNotEmpty=/service/

[Install]
WantedBy=multi-user.target
