# source/systemd/service-manager-svscan.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Daemontools compatible service scanner
After=service-manager.socket
After=service-manager.service

[Service]
ExecStart=/bin/service-dt-scanner /service/
