# source/systemd/system/service-manager-svscan.path
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Daemontools compatible /service/ directory

[Path]
DirectoryNotEmpty=/service/

[Install]
WantedBy=multi-user.target
