# source/systemd/user/service-manager.socket
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Service manager local domain socket activation
#After=systemd-tmpfiles-setup.service

[Socket]
ListenDatagram=%t/%p/control
Accept=no
SocketMode=0700

[Install]
WantedBy=sockets.target
