# source/anacron.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=anacron scheduling service

[Service]
systemdWorkingDirectory=false
ExecStartPre=on_ac_power
ExecStart=%p -d -s
Restart=always

[Install]
WantedBy=workstation.target
