# source/atd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description="at" scheduling service
[Install]
WantedBy=workstation.target
[Service]
systemdWorkingDirectory=false
ExecStart=atd -f
Restart=always
