# source/cups.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Common Unix Printing Service

[Service]
systemdWorkingDirectory=false
ExecStart=cupsd -f
Restart=always

[Install]
WantedBy=workstation.target
