# source/dmesg.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Dump kernel boot log sysinit service
DefaultDependencies=false
EarlySupervise=true

[Service]
Umask=0022
systemdWorkingDirectory=false
ExecStart=fdredir --write 1 /run/dmesg.boot dmesg
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
