# source/epmd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Erlang Port Mapper service

[Service]
systemdWorkingDirectory=false
User=rabbitmq
ExecStart=%p
Restart=always

[Install]
WantedBy=workstation.target
