# source/fsck-root.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Filesystem check the / volume
DefaultDependencies=false
EarlySupervise=true
Before=remount-root.service

[Service]
Type=oneshot
ExecStart=fsck -C -p /
systemdWorkingDirectory=false

[Install]
WantedBy=remount-root.service
