# source/fsck-usr.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Filesystem check the /usr volume
DefaultDependencies=false
EarlySupervise=true
Before=mount-usr.service

[Service]
Type=oneshot
ExecStart=fsck -C -p /usr
systemdWorkingDirectory=false

[Install]
WantedBy=mount-usr.service
