# source/ftp4d.socket
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]

[Socket]
Backlog=2
MaxConnections=2
Nagle=false
ListenStream=0:21
UCSPIRules=true
Accept=true

[Install]
WantedBy=server.target
