# source/ipxrouted.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=IPX route daemon

[Service]
systemdWorkingDirectory=false
ExecStart=IPXrouted
Restart=always

[Install]
WantedBy=server.target
