# source/kerneloops.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Tool to automatically collect and submit crash signatures
[Install]
WantedBy=workstation.target
[Service]
systemdWorkingDirectory=false
ExecStart=kerneloops --nodaemon
Restart=always
