# source/kmod@unix.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=kmod loader for unix

[Service]
Type=oneshot
systemdWorkingDirectory=false
ExecStart=[ -d /proc/net/unix/ ] || kmod modprobe unix
Restart=always

[Install]
WantedBy=workstation.target
