# source/kmod@unix.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=kmod loader for unix

[Service]
systemdWorkingDirectory=false
ExecStart=sh -c '[ -d /proc/net/unix/ ] || kmod modprobe unix'
RemainAfterExit=yes

[Install]
WantedBy=workstation.target
