# source/kmod@autofs.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=kmod loader for autofs4

[Service]
Type=oneshot
systemdWorkingDirectory=false
ExecStart=sh -c '[ -d /sys/class/misc/autofs4/ ] || kmod modprobe autofs4'
RemainAfterExit=yes

[Install]
WantedBy=workstation.target
