# source/mroute6d.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=IPv6 mroute daemon

[Service]
systemdWorkingDirectory=false
ExecStart=mroute6d
Restart=always

[Install]
WantedBy=server.target
