# source/name-services.target
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=DNS milestone
DefaultDependencies=false

[Service]
Type=oneshot
User=nobody
systemdWorkingDirectory=false
RemainAfterExit=yes

[Install]
