# source/ntpd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=BSD NTP daemon
[Install]
WantedBy=workstation.target
[Service]
systemdWorkingDirectory=false
ExecStart=ntpd -n -g -f /var/db/ntpd.drift
Restart=always
