# source/rarpd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=RARP service

[Service]
systemdWorkingDirectory=false
ExecStart=%p -d -f -a
Restart=always

[Install]
WantedBy=server.target
