# source/saned@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Scanner service

[Service]
Group=scanner
systemdWorkingDirectory=false
ExecStart=%P
StandardInput=socket
Restart=always

[Install]
WantedBy=server.target
