# source/service-manager-svscan.service
[Unit]
Description=Daemontools compatible service scanner

[Service]
ExecStart=/bin/service-dt-scanner /service/
