# source/service-manager.service
[Unit]
Description=Service manager

[Service]
ExecStart=/bin/service-manager
TimeoutSec=10
