# source/service-manager-svscan.path
[Path]
DirectoryNotEmpty=/service/
After=local-fs.target sysinit.target

[Install]
WantedBy=multi-user.target
