# source/sickbeard.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=SickBeard daemon

[Service]
User=sickbeard
WorkingDirectory=/opt/sickbeard
ExecStart=sh -c 'exec python SickBeard.py --nolaunch ${SB_OPTS}'

[Install]
WantedBy=workstation.target
