# source/sshd@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=SSH per-connection service
#Requires=sshdgenkeys.service
After=sshdgenkeys.service

[Service]
systemdWorkingDirectory=false
ExecStart=%P -i -e
Restart=always

[Install]
WantedBy=server.target
