# source/sshd@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=SSH per-connection service
Requires=sshdgenkeys.service
After=sshdgenkeys.service
[Install]
WantedBy=server.target
[Service]
systemdWorkingDirectory=false
ExecStart=sshd -i -e
Restart=always
