# source/swapauto.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=%P service

[Service]
systemdWorkingDirectory=false
ExecStartPre=swapon -a
ExecStart=true
RemainAfterExit=yes

[Install]
WantedBy=workstation.target
