# source/swapexd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=BSD swapex daemon

[Service]
systemdWorkingDirectory=false
ExecStart=swapexd -p
Restart=always

[Install]
WantedBy=workstation.target
