# source/sysinit.target
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Sysinit group target
DefaultDependencies=false
EarlySupervise=true
Before=shutdown.target

[Service]
User=nobody
systemdWorkingDirectory=false
ExecStart=pause
Restart=always

[Install]
StoppedBy=shutdown.target
