# source/syslogd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Old-style system log daemon

[Service]
systemdWorkingDirectory=false
ExecStart=syslogd -c -ss
Restart=always

[Install]
WantedBy=workstation.target
