# source/uhidd@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=%P service for device %I

[Service]
systemdWorkingDirectory=false
ExecStart=%P %I
Restart=always

[Install]
WantedBy=workstation.target
