# source/uuidd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=util-linux UUID listening daemon

[Service]
systemdWorkingDirectory=false
systemdUserEnvironment=no
User=uuidd
ExecStartPre=mkdir -m 0755 /run/%p
ExecStartPre=chown %p:%p /run/%p
ExecStart=%p --socket-activation

[Install]
WantedBy=workstation.target
