# source/VBoxService.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=VirtualBox guest additions daemon

[Service]
systemdWorkingDirectory=false
ExecStart=VBoxService -f
Restart=always

[Install]
WantedBy=virtualbox.target
