# source/volmand.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=BSD volume manager daemon
[Install]
WantedBy=workstation.target
[Service]
systemdWorkingDirectory=false
ProcessGroupLeader=true
ExecStart=volmand daemon
Restart=always
