# source/watchdogd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=BSD kernel watchdog daemon

[Service]
systemdWorkingDirectory=false
ExecStart=%p -d
Restart=always

[Install]
WantedBy=workstation.target
