# source/winbindd.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Name service switch daemon for SAMBA

[Service]
systemdWorkingDirectory=false
ExecStart=%p -F -s /usr/local/etc/smb.conf
Restart=always

[Install]
WantedBy=server.target
