# source/wpa_supplicant.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Wi-Fi Protected Access client daemon
After=dbus.service

[Service]
systemdWorkingDirectory=false
ExecStart=%P -u
Restart=always

[Install]
WantedBy=workstation.target
