# source/convert/ataidle@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Apply %p settings to disc %I
Wants=%p-log.service
After=%p-log.service
DefaultDependencies=no

[Service]
Type=oneshot
EnvironmentDirectory=env
RemainAfterExit=true
ExecStart=%p ${args} /dev/"%I"

[Install]
WantedBy=workstation.target
