# source/convert/dirmngr.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=GnuPG network certificate management daemon

[Service]
RuntimeDirectory=gnupg
ExecStart=%p --supervised

[Install]
