# source/convert/intrat.target
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Intrat target
Before=shutdown

[Install]
StoppedBy=shutdown
