# source/convert/knotify.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=KDE Notify service for %U

[Service]
Type=dbus
EnvironmentDirectory=env
ExecStart=%p4
Restart=on-abort

[Install]
