# source/convert/socket-servers-log.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Per-user socket servers' log service for %U

[Service]
WorkingDirectory=/var/log/user/
ExecStart=cyclog "%U/socket-servers/"

[Install]
