# source/convert/spppcontrol@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=BSD %p for %i
Wants=%p-log.service
After=%p-log.service

[Service]
Type=oneshot
EnvironmentDirectory=env
ExecStart=%p "%I" ${args}

[Install]
WantedBy=static-networking.target
