# source/convert/spppcontrol@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=BSD %p for %i

[Service]
Type=oneshot
EnvironmentDirectory=env
ExecStart=%p "%I" ${args}

[Install]
WantedBy=workstation.target
