# source/convert/static_arp@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=BSD %p for %i

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentDirectory=env
ExecStart=ndp -s %I
ExecStopPost=set -- %i && ndp -d $1

[Install]
WantedBy=server.target
