# source/convert/static_ndp@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=BSD static NDP binding for %i

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentDirectory=env
ExecStart=route -6 add ${route}
ExecStopPost=route -6 del ${route}
