# source/convert/tinydns.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Bernstein's content DNS server daemon, modified to use systemd LISTEN_FDS protocol
After=network-interfaces.service
Before=name-services.target

[Service]
User=tinydns-d
EnvironmentUserOnly=true
EnvironmentDirectory=env
LimitDATA=3000000
ExecStart=%p
RestartSec=0.2

[Install]
WantedBy=workstation.target
