# source/convert/user@.target
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Per-user %I target

[Service]
User=%I

[Install]
# A per-user target is not auto-startable at bootstrap.
