# source/convert/user-dbus@.socket
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Per-user "Desktop" bus socket
After=user-runtime@%I
Wants=user-runtime@%I

[Socket]
Backlog=2
MaxConnections=2
ListenStream=/run/user/%I/dbus/user_bus_socket
Accept=false

[Install]
WantedBy=sockets.target
