# source/convert/user-runtime@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Per-user runtime for %I

[Service]
RemainAfterExit=true
User=%I

[Install]
WantedBy=user@%I.target
