# source/convert/user@.target
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Per-user %I target

[Service]
User=%I
# This is a pseudo-daemon that does not have a running process.
RemainAfterExit=yes
Restart=always

[Install]
# A per-user target is not auto-startable at bootstrap.
