# source/convert/warden-jailed@.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Run command within PC-BSD Warden jail %I

[Service]
JailID=%I
ExecStart=true

[Install]
WantedBy=warden.target
